# storage/mroonga/vendor/groonga/data/systemd/fedora/groonga-server-gqtp.service
[Unit]
Description=Groonga Text Searching Engine
After=network.target

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/groonga-server-gqtp
User=groonga
Group=groonga
ExecStart=/usr/bin/groonga -d --pid-path /var/run/groonga/groonga-gqtp.pid --bind-address $ADDRESS --log-path $LOG_PATH --query-log-path $QUERY_LOG_PATH --protocol $PROTOCOL --port $PORT $DATABASE

[Install]
WantedBy=multi-user.target
